# roles/services/files/docker-cleanup.timer
[Unit]
Description=Remove unused docker images weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
